# supervisor/conf.d/mysql.conf
[program:mysql]

command=/usr/local/bin/mysqld

user=sabo
stdout_logfile=/data/log/mysql/log
redirect_stderr=true
