# supervisor/hiveserver2.conf
[program:hiveserver2]

command=/Users/sabo/Developer/hive-0.12.0/bin/hiveserver2

user=sabo
stdout_logfile=/data/log/hiveserver2/log
redirect_stderr=true
